# systemd/user/wall_change.service
[Unit]
Description=Set and change wallpaper

[Service]
ExecStart=/bin/sh /home/kouta/Scripts/wall_change.sh
